# includes/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

[program:displaylogs]
command=tail -f /supervisord.log /var/log/nginx/passgen_access.log /var/log/nginx/passgen_error.log